# miner.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/miner.sh
User=root
Group=root

[Install]
WantedBy=default.target
